# extra/init.d/pierre
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net nqptp avahi-daemon
	use logger 
}

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start \
    --user thughey:thughey \
    --pidfile /run/pierre/${SVCNAME}.pid \
		--exec /home/thughey/devel/pierre/build/${SVCNAME} -- \
	  --daemon
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --pidfile /run/pierre/${SVCNAME}.pid \
		--retry SIGINT/5 SIGTERM/5
	eend $?
}
